# pkg/arch/cts.install
post_install() {
    systemctl enable cts.service
}
