# cts.service
[Unit]
Description=Chronos server daemon

[Service]
ExecStart=/usr/bin/cts

[Install]
WantedBy=multi-user.target
